# RaspberryPi/NetworkConfig/rpi_wifi_watchdog.service
[Unit]
Description=Raspberry Pi Wi-Fi Watchdog (Hotspot + Fallback)
After=network.target NetworkManager.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rpi_wifi_watchdog.sh
